# scripts/maintenance/systemd/lab-docker-container-cleanup.timer
[Unit]
Description=Daily laboratory Docker container cleanup

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=15m
Persistent=true
Unit=lab-docker-container-cleanup.service

[Install]
WantedBy=timers.target
